# supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=*:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
;serverurl=unix://tmp/supervisor.sock
serverurl=http://127.0.0.1:9001

;[program:hbase]
;environment=JAVA_HOME=/usr
;command=/opt/hbase/bin/start-hbase.sh

[program:hbase-master]
environment=JAVA_HOME=/usr
command=/opt/hbase/bin/hbase master start

[program:hbase-regionserver]
command=/opt/hbase/bin/hbase regionserver start
